# src/services/main/apache.1.service
[Unit]
Description=Apache instance 1
After=docker.service
Requires=docker.service apache.1.sk.service
Before=apache.1.sk.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=0
ExecStartPre=-/opt/docker/bin/docker kill apache.1
ExecStartPre=-/opt/docker/bin/docker rm apache.1
ExecStart=/opt/docker/bin/docker run --name apache.1 -p 8085:80 local/apache2
ExecStop=/opt/docker/bin/docker stop -t 20 apache.1
