# src/services/main/apache@.service
[Unit]
Description=Apache Service Instance
After=docker.service
Requires=docker.service apache.sk@%i.service apache.hc@%i.timer
Before=apache.sk@%i.service apache.hc@%i.timer

[Service]
TimeoutStartSec=0
ExecStartPre=-/opt/docker/bin/docker kill apache.%i
ExecStartPre=-/opt/docker/bin/docker rm apache.%i
ExecStart=/opt/docker/util/dockerlaunch apache.%i local/apache2
ExecStop=/opt/docker/bin/docker stop -t 20 apache.%i
Restart=always
RestartSec=4
