# src/services/main/apache.endpoint.service
[Unit]
Description=Announce Apache Endpoint

[Service]
ExecStart=/opt/docker/util/announce.endpoint
ExecStop=/opt/docker/util/unannounce.endpoint
