# src/services/main/apache@.service
[Unit]
Description=Apache Instance %i
After=docker.service
Requires=docker.service apache.sk@%i.service
Before=apache.sk@%i.service

[Service]
EnvironmentFile=/opt/docker/conf/apache.%i.sh
TimeoutStartSec=0
ExecStartPre=-/opt/docker/bin/docker kill apache.%i
ExecStartPre=-/opt/docker/bin/docker rm apache.%i
ExecStart=/bin/sh -c '/opt/docker/bin/docker run --name apache.%i -p $PORT:80 local/apache2'
ExecStop=/opt/docker/bin/docker stop -t 20 apache.%i
