# src/services/main/apache.2.service
[Unit]
Description=Apache instance 2
After=docker.service
Requires=docker.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=0
ExecStartPre=-/opt/docker/bin/docker kill apache.2
ExecStartPre=-/opt/docker/bin/docker rm apache.2
ExecStart=/opt/docker/bin/docker run --name apache.2 -p 8086:80 local/apache2
ExecStop=/opt/docker/bin/docker stop -t 20 apache.2
