# src/services/sidekicks/apache.2.sk.service
[Unit]
Description=Announce Apache2
BindsTo=apache.2.service
After=apache.2.service

[Service]
EnvironmentFile=/etc/environment
TimeoutStartSec=0
ExecStart=/opt/docker/util/announce.vulcan 2 8086
ExecStop=/opt/docker/util/unannounce.vulcan 2

[X-Fleet]
X-ConditionMachineOf=apache.2.service
