# src/services/healthchecks/apache.hc@.service
[Unit]
Description=Check the health of the Apache Service
BindsTo=apache@%i.service
After=apache@%i.service

[Service]
EnvironmentFile=/opt/docker/conf/apache.%i.sh
TimeoutStartSec=0
ExecStart=-/bin/sh -c '/opt/docker/util/healthcheck.http $PORT'

[X-Fleet]
X-ConditionMachineOf=apache@%i.service
