# src/base/docker.service
[Unit]
Description=Docker Service

[Service]
ExecStart=/opt/docker/bin/docker -d
