# src/base/fleetweb.service
[Unit]
Description=Fleet Read Only Web Interface
After=network.target fleetd.service

[Service]
ExecStart=/bin/bash -c "cd /opt/docker/src/fleetweb/ && /opt/docker/bin/fleetweb"

[Install]
RequiredBy=multi-user.target
